# config/etc/supervisor/conf.d/etcd.conf
[program:etcd]
command=/usr/local/bin/etcd
redirect_stderr=true
stopsignal=TERM
autorestart=true
